# etc/supervisor/conf.d/exim4.conf
[program:exim4]
command=/usr/sbin/exim4 -bdf -q30m
priority=1
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

